# cli/packaging/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create the config directory if it does not exist yet.
        mkdir -p /etc/openme

        # Reload systemd and enable the service unit introduced by this package.
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload || true
            systemctl enable openme.service || true
            echo "openme: service installed. Start it with:"
            echo "  sudo openme init --server <hostname>"
            echo "  sudo systemctl start openme"
        fi
        ;;
esac

exit 0
